# .config/systemd/user/wallpaper-changer.timer
[Unit]
Description=Run wallpaper-changer every day

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target